# chromeos-base/acpi-patch-service/files/acpi-patch-service.conf
# Copyright 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Add custome acpi patch for board"
author          "yang@fydeos.io"

# This job is detecting the specific usb device to determine if is tablet mode
# enable.
#

start on starting boot-services

oom score -100

exec /usr/share/fydeos_shell/patch_acpi.sh
